# util/conf/docker.service.d/vagrant.conf
[Unit]
After=vagrant.mount
Wants=vagrant.mount
